# contrib/init/bitcoinpurpled.conf
description "BitcoinPurple Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BITCOINPURPLED_BIN="/usr/bin/bitcoinpurpled"
env BITCOINPURPLED_USER="bitcoinpurple"
env BITCOINPURPLED_GROUP="bitcoinpurple"
env BITCOINPURPLED_PIDDIR="/var/run/bitcoinpurpled"
# upstart can't handle variables constructed with other variables
env BITCOINPURPLED_PIDFILE="/var/run/bitcoinpurpled/bitcoinpurpled.pid"
env BITCOINPURPLED_CONFIGFILE="/etc/bitcoinpurple/bitcoinpurple.conf"
env BITCOINPURPLED_DATADIR="/var/lib/bitcoinpurpled"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # bitcoinpurpled will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BITCOINPURPLED_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run bitcoinpurpled."
        echo "The setting must appear in $BITCOINPURPLED_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"BitcoinPurple Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BITCOINPURPLED_PIDDIR"
    chmod 0755 "$BITCOINPURPLED_PIDDIR"
    chown $BITCOINPURPLED_USER:$BITCOINPURPLED_GROUP "$BITCOINPURPLED_PIDDIR"
    chown $BITCOINPURPLED_USER:$BITCOINPURPLED_GROUP "$BITCOINPURPLED_CONFIGFILE"
    chmod 0660 "$BITCOINPURPLED_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BITCOINPURPLED_PIDFILE" \
    --chuid $BITCOINPURPLED_USER:$BITCOINPURPLED_GROUP \
    --exec "$BITCOINPURPLED_BIN" \
    -- \
    -pid="$BITCOINPURPLED_PIDFILE" \
    -conf="$BITCOINPURPLED_CONFIGFILE" \
    -datadir="$BITCOINPURPLED_DATADIR" \
    -disablewallet \
    -daemon

